# blog_boluo.conf
[program:blog_boluo]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py --bind 0.0.0.0:2000 --pid /tmp/blog_boluo.pid
directory=/var/www/blog_boluo
autostart=true
autorestart=true